# twimg2rss.crontab
LANG=ja_JP.UTF-8
2,7,12,17,22,27,32,37,42,47,52,57 * * * * INSTALL_USER /opt/twimg2rss.py >> LOGDIR/cron.log 2>&1
0 * * * * INSTALL_USER find TIMELINE -name 'timeline_[0-9-]*.json' -mtime +7 -exec rm {} \;
